# infra/supervisor/customerflow.conf
; CustomerFlow AI — Supervisor configuration (Contabo VPS + aaPanel).
;
; Drop into /etc/supervisor/conf.d/customerflow.conf, then:
;   supervisorctl reread && supervisorctl update
;
; Each program logs to /var/log/customerflow/<name>.{out,err}.log.
; Configure log rotation in /etc/logrotate.d/customerflow if you don't already
; have a global rule.

[group:customerflow]
programs=customerflow-api,customerflow-worker,customerflow-web,customerflow-loyalty-pwa
priority=10

; ── FastAPI (Uvicorn behind Gunicorn) ───────────────────────────────────────
[program:customerflow-api]
command=/opt/customerflow/.venv/bin/gunicorn app.main:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --graceful-timeout 30 \
    --max-requests 5000 \
    --max-requests-jitter 500 \
    --access-logfile - \
    --error-logfile -
directory=/opt/customerflow/apps/api
user=customerflow
group=customerflow
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1",ENVIRONMENT="production"
stdout_logfile=/var/log/customerflow/api.out.log
stderr_logfile=/var/log/customerflow/api.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; ── ARQ Worker (background tasks) ───────────────────────────────────────────
[program:customerflow-worker]
command=/opt/customerflow/.venv/bin/arq app.workers.worker_settings.WorkerSettings
directory=/opt/customerflow/apps/api
user=customerflow
group=customerflow
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1",ENVIRONMENT="production"
stdout_logfile=/var/log/customerflow/worker.out.log
stderr_logfile=/var/log/customerflow/worker.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; ── Next.js (production server) ─────────────────────────────────────────────
[program:customerflow-web]
command=/usr/bin/node /opt/customerflow/apps/web/.next/standalone/server.js
directory=/opt/customerflow/apps/web
user=customerflow
group=customerflow
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1",NEXT_TELEMETRY_DISABLED="1",LOYALTY_PWA_ORIGIN="http://127.0.0.1:3001",API_PROXY_TARGET="http://127.0.0.1:8000"
stdout_logfile=/var/log/customerflow/web.out.log
stderr_logfile=/var/log/customerflow/web.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; ── Loyalty PWA (customer rewards wallet at /rewards/*) ───────────────────────
[program:customerflow-loyalty-pwa]
command=/usr/bin/node node_modules/next/dist/bin/next start -p 3001
directory=/opt/customerflow/apps/loyalty-pwa
user=customerflow
group=customerflow
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
environment=NODE_ENV="production",PORT="3001",HOSTNAME="127.0.0.1",NEXT_TELEMETRY_DISABLED="1",API_PROXY_TARGET="http://127.0.0.1:8000",NEXT_PUBLIC_API_URL=""
stdout_logfile=/var/log/customerflow/loyalty-pwa.out.log
stderr_logfile=/var/log/customerflow/loyalty-pwa.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; Optional: a periodic cleanup job (in case you don't want to run ARQ cron).
; [program:customerflow-cron]
; command=/bin/bash -c "while true; do /opt/customerflow/.venv/bin/python -m app.scripts.cleanup; sleep 3600; done"
; user=customerflow
; autostart=true
; autorestart=true
